# dot_config/systemd/user/backup-irclogs.timer
[Unit]
Description=Backup irc logs

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
