# dot_config/systemd/user/backup-irclogs.service
[Unit]
Description=backup IRC logs
Wants=network-online.target
After=network-online.target
ConditionFileNotEmpty=/media/t/Archives/irc-logs/backup.sh
Before=borgmatic@tosh.service

[Install]
WantedBy=borgmatic@tosh.service

[Service]
Type=oneshot
LockPersonality=true
MemoryDenyWriteExecute=no
NoNewPrivileges=yes
PrivateTmp=yes
ProtectControlGroups=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ProtectSystem=full

Nice=19
CPUSchedulingPolicy=batch
IOSchedulingClass=best-effort
IOSchedulingPriority=7
IOWeight=100

Restart=no
LogRateLimitIntervalSec=0

ExecStart=/media/t/Archives/irc-logs/backup.sh
